# systemd/aitbc-openclaw.service.d/10-central-env.conf
[Service]
EnvironmentFile=/opt/aitbc/.env
