# apps/ffmpeg/aitbc-ffmpeg.service
[Unit]
Description=AITBC FFmpeg Video Processing Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
Environment=FFMPEG_PORT=8230
Environment=FFMPEG_GPU_DEVICE=0
Environment=FFMPEG_HW_ACCEL=cuda
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/ffmpeg/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-ffmpeg

[Install]
WantedBy=multi-user.target
