# apps/miner-dashboard/aitbc-miner-dashboard.service
[Unit]
Description=AITBC Miner Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc-miner-dashboard
Environment=PYTHONPATH=/opt/aitbc-miner-dashboard
ExecStart=/opt/aitbc-miner-dashboard/.venv/bin/python dashboard_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
