# apps/miner/aitbc-miner.service
[Unit]
Description=AITBC Production GPU Miner
After=network.target ollama.service aitbc-coordinator-api.service

[Service]
Type=simple
User=aitbc
WorkingDirectory=/opt/aitbc/apps/miner
Environment="PATH=/opt/aitbc/venv/bin:/usr/bin:/usr/local/bin"
Environment="PYTHONPATH=/opt/aitbc"
Environment="COORDINATOR_URL=http://localhost:8203"
Environment="MINER_ID=aitbc-miner-1"
Environment="LOG_DIR=/var/log/aitbc"
# Secret: AUTH_TOKEN moved to /etc/aitbc/%N.env (created at deploy time, not committed)
EnvironmentFile=/etc/aitbc/%N.env
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/miner/production_miner.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-miner
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Security hardening (v0.5.0)
ProtectSystem=full
ReadWritePaths=/opt/aitbc /var/lib/aitbc /var/log/aitbc /run/aitbc

# Security hardening (v0.5.0)
# WatchdogSec=30 # Disabled - requires application-level sd_notify support

[Install]
WantedBy=multi-user.target
