# apps/exchange/aitbc-exchange.service
[Unit]
Description=AITBC Exchange API Service
After=network.target aitbc-blockchain-node.service
Wants=aitbc-blockchain-node.service

[Service]
Type=simple
User=aitbc-internal
Group=aitbc-services
WorkingDirectory=/opt/aitbc/apps/exchange
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
Environment="PYTHONPATH=/opt/aitbc"
ExecStart=/opt/aitbc/venv/bin/python simple_exchange_api.py --port 8106
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-exchange-api

# Note: ProtectSystem=no removed for service user compatibility

[Install]
WantedBy=multi-user.target
