# apps/exchange/aitbc-exchange.service
[Unit]
Description=AITBC Exchange API Service
After=network.target aitbc-blockchain-node.service
Wants=aitbc-blockchain-node.service

[Service]
Type=simple
User=aitbc-internal
Group=aitbc-services
WorkingDirectory=/opt/aitbc/apps/exchange
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
Environment="PYTHONPATH=/opt/aitbc:/opt/aitbc/aitbc"
Environment="BLOCKCHAIN_RPC_URL=http://localhost:8202"
# Config: REDIS_URL moved to /etc/aitbc/%N.env (created at deploy time, not committed)
EnvironmentFile=/etc/aitbc/%N.env
# Network hardening (v0.5.1): Exchange service binds to localhost via HTTPServer
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/exchange/simple_exchange_api.py --port 8106
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-exchange-api
# Note: ProtectSystem=no removed for service user compatibility
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Network hardening (v0.5.1)
IPDeny=any
# Security hardening (v0.5.0)
ProtectSystem=full
ReadWritePaths=/opt/aitbc /var/lib/aitbc /var/log/aitbc /run/aitbc

# Security hardening (v0.5.0)
# WatchdogSec=30 # Disabled - requires application-level sd_notify support
