# scripts/utils/aitbc-load-secrets.service
[Unit]
Description=AITBC Load Secrets Service
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
ExecStart=/opt/aitbc/scripts/utils/load-keystore-secrets.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-load-secrets
RemainAfterExit=yes
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
ProtectSystem=full
ReadWritePaths=/run/aitbc /var/log/aitbc
