# scripts/utils/aitbc-plugin.service
[Unit]
Description=AITBC Plugin Service
After=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc
Environment=PATH=/opt/aitbc/venv/bin:/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment="PYTHONPATH=/opt/aitbc/apps/agent-management/examples/plugin-service/src"
Environment="DATA_DIR=/var/lib/aitbc"
Environment="LOG_DIR=/var/log/aitbc"
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
ExecStart=/opt/aitbc/venv/bin/python -m uvicorn plugin_service.main:app --host 0.0.0.0 --port 8109 --log-level warning
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-plugin
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Security hardening (v0.5.0)
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
