# scripts/systemd/aitbc-recovery.service
[Unit]
Description=AITBC Service Recovery — relinks systemd units and loads secrets after reboot
After=network.target
Before=aitbc-blockchain-node.service aitbc-coordinator-api.service aitbc-marketplace.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/opt/aitbc/scripts/utils/link-systemd.sh && /opt/aitbc/scripts/utils/load-keystore-secrets.sh'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Security hardening (v0.5.0)
ProtectSystem=full
ReadWritePaths=/opt/aitbc /var/lib/aitbc /var/log/aitbc /run/aitbc
# Security hardening (v0.5.0)
RestartSec=5
# Security hardening (v0.5.0)
Restart=on-failure

# Security hardening (v0.5.0)
WatchdogSec=30
NotifyAccess=all

[Install]
WantedBy=multi-user.target
