# scripts/systemd/aitbc-marketplace.service
[Unit]
Description=AITBC Marketplace Service (Port 8102)
After=network.target aitbc-blockchain-rpc.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/aitbc:/opt/aitbc/apps/marketplace/src
Environment=MARKETPLACE_DATABASE_URL=sqlite+aiosqlite:////var/lib/aitbc/data/marketplace_service.db
ExecStart=/opt/aitbc/venv/bin/python -m marketplace_service.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-marketplace

# Allow database writes for SQLite WAL mode
ProtectSystem=no

[Install]
WantedBy=multi-user.target
