# scripts/systemd/aitbc-recovery.service
[Unit]
Description=AITBC Service Recovery — relinks systemd units and loads secrets after reboot
After=network.target
Before=aitbc-blockchain-node.service aitbc-coordinator-api.service aitbc-marketplace.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c '/opt/aitbc/scripts/utils/link-systemd.sh && /opt/aitbc/scripts/utils/load-keystore-secrets.sh'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
