# apps/agent-coordinator/aitbc-agent-coordinator.service
[Unit]
Description=AITBC Agent Coordinator Service
After=network.target redis.service aitbc-blockchain-node.service
Wants=aitbc-blockchain-node.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/aitbc:/opt/aitbc/apps/agent-coordinator/src
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
EnvironmentFile=/etc/aitbc/blockchain-secrets.env
# Agent coordinator execution
ExecStart=/opt/aitbc/venv/bin/python -m app.main
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10
# Production reliability
Restart=on-failure
RestartSec=5
# Production logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-agent-coordinator
# Production security
NoNewPrivileges=yes
# Production performance
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=256M
CPUQuota=50%
# Security hardening (v0.5.0)
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Security hardening (v0.5.0)
ProtectSystem=full
ReadWritePaths=/opt/aitbc /var/lib/aitbc /var/log/aitbc /run/aitbc

# Security hardening (v0.5.0)
# WatchdogSec=30 # Disabled - requires application-level sd_notify support

[Install]
WantedBy=multi-user.target
