# apps/whisper/aitbc-whisper.service
[Unit]
Description=AITBC Whisper Transcription Service
After=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/whisper
Environment="PATH=/opt/aitbc/venv/bin"
Environment="PYTHONPATH=/opt/aitbc/packages/py/aitbc-core/src:/opt/aitbc/apps/whisper:$PYTHONPATH"
Environment="WHISPER_MODEL=base"
Environment="WHISPER_DEVICE=cuda"
Environment="WHISPER_COMPUTE_TYPE=float16"
Environment="HF_HOME=/var/lib/aitbc/whisper-cache"
Environment="TRANSFORMERS_CACHE=/var/lib/aitbc/whisper-cache"
ExecStart=/opt/aitbc/venv/bin/python main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-whisper
MemoryMax=1G
# Security hardening (v0.5.0)
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
# Security hardening (v0.5.0)
ProtectSystem=full
ReadWritePaths=/opt/aitbc /var/lib/aitbc /var/log/aitbc /run/aitbc

# Security hardening (v0.5.0)
# WatchdogSec=30 # Disabled - requires application-level sd_notify support
