# systemd/aitbc-marketplace.service
[Unit]
Description=AITBC Agent-First GPU Marketplace
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc/apps/marketplace
Environment=PYTHONPATH=/opt/aitbc/apps/marketplace:/opt/aitbc/packages/py/aitbc-sdk/src
ExecStart=/usr/bin/python3 /opt/aitbc/apps/marketplace/agent_marketplace.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
