# systemd/aitbc-host-gpu-miner.service
[Unit]
Description=AITBC Host GPU Miner Client
After=network.target
Wants=network.target

[Service]
Type=simple
User=oib
Group=oib
WorkingDirectory=/opt/gpu-miner
Environment=PYTHONPATH=/opt/gpu-miner
Environment=VIRTUAL_ENV=/opt/gpu-miner/.venv
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8
Environment=FORCE_COLOR=1
Environment=PATH=/opt/gpu-miner/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -c 'exec /opt/gpu-miner/.venv/bin/python -u /opt/gpu-miner/scripts/gpu/gpu_miner_host.py 2>&1 | tee -a /opt/gpu-miner/logs/host_gpu_miner.log'
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-host-gpu-miner

[Install]
WantedBy=multi-user.target
