# systemd/aitbc-exchange.service
[Unit]
Description=AITBC Trade Exchange
After=network.target

[Service]
Type=exec
User=debian
WorkingDirectory=/opt/aitbc/apps/exchange
Environment=PATH=/usr/bin
ExecStart=/usr/bin/python3 server.py --port 3002
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
