# systemd/aitbc-coordinator-proxy-health.service
[Unit]
Description=AITBC Coordinator Proxy Health Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/aitbc/apps/coordinator-api/scripts/check_coordinator_proxy.sh

[Install]
WantedBy=multi-user.target
