# systemd/aitbc-multimodal.service
[Unit]
Description=AITBC Multi-Modal Agent Service
After=network.target aitbc-coordinator-api.service
Wants=aitbc-coordinator-api.service

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/opt/aitbc/apps/coordinator-api
Environment=PATH=/opt/aitbc/venv/bin:/usr/bin
Environment=PYTHONPATH=/opt/aitbc/apps/coordinator-api/src
EnvironmentFile=/etc/aitbc/.env
ExecStart=/opt/aitbc/venv/bin/python -m uvicorn src.app.services.multimodal_app:app --host 127.0.0.1 --port 8020
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-multimodal

# Security
NoNewPrivileges=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
