# systemd/aitbc-blockchain-sync.service
[Unit]
Description=AITBC Blockchain Chain Synchronization Service
After=network.target redis.service aitbc-blockchain-node.service aitbc-blockchain-rpc.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
EnvironmentFile=/etc/aitbc/.env
EnvironmentFile=/etc/aitbc/node.env
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/scripts/wrappers/aitbc-blockchain-sync-wrapper.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
