# systemd/aitbc-web-ui.service
[Unit]
Description=AITBC Web UI Service (Port 8016)
Documentation=https://docs.aitbc.bubuit.net
After=network.target aitbc-coordinator-api.service
Wants=aitbc-coordinator-api.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/explorer-web/dist
Environment=PATH=/usr/bin:/bin
Environment=PYTHONPATH=/opt/aitbc/apps/explorer-web/dist
Environment=PORT=8016
Environment=SERVICE_TYPE=web-ui
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 -m http.server 8016 --bind 127.0.0.1
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-web-ui

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/aitbc /var/lib/aitbc/data
LimitNOFILE=65536

# Resource limits
MemoryMax=1G
CPUQuota=100%

[Install]
WantedBy=multi-user.target
