# systemd/aitbc-monitor.service
[Unit]
Description=AITBC Monitor Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/aitbc/.env
EnvironmentFile=/etc/aitbc/node.env

# Monitor execution
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/monitor/monitor.py

# Production reliability
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Production logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-monitor

# Production security
NoNewPrivileges=true
ProtectHome=true

# Production performance
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
