# systemd/aitbc-agent-coordinator.service
[Unit]
Description=AITBC Agent Coordinator Service
After=network.target aitbc-agent-registry.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc/apps/agent-services/agent-coordinator/src
Environment=PYTHONPATH=/opt/aitbc
ExecStart=/opt/aitbc/venv/bin/python coordinator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
