# systemd/aitbc-blockchain-sync-dev.service
[Unit]
Description=AITBC Blockchain Chain Synchronization Service (Dev)
After=network.target redis.service aitbc-blockchain-node-dev.service aitbc-blockchain-rpc-dev.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc/apps/blockchain-node
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/aitbc/apps/blockchain-node/src:/opt/aitbc/apps/blockchain-node/scripts
EnvironmentFile=/etc/aitbc/blockchain.env
ExecStart=/usr/bin/python3 -m aitbc_chain.chain_sync \
  --redis redis://localhost:6379 \
  --node-id follower-dev \
  --rpc-port 8026 \
  --leader-host 127.0.0.1 \
  --source-host 127.0.0.1 --source-port 8006 \
  --import-host 127.0.0.1 --import-port 8026
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
