# systemd/aitbc-gpu-miner.service
[Unit]
Description=AITBC GPU Miner Client
After=network.target
Wants=network.target

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/opt/gpu-miner
Environment=PYTHONPATH=/opt/gpu-miner
ExecStart=/opt/gpu-miner/venv/bin/python /opt/gpu-miner/gpu_miner_simple.py
Restart=on-failure
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-gpu-miner

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/gpu-miner

[Install]
WantedBy=multi-user.target
