# systemd/aitbc-marketplace.service
[Unit]
Description=AITBC Marketplace Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/marketplace-service
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/root/.local/bin/poetry run python -m marketplace_service.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-marketplace

[Install]
WantedBy=multi-user.target
