# systemd/aitbc-exchange-mock-api.service
[Unit]
Description=AITBC Exchange Mock API (trades/orderbook/wallet-connect)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /var/www/aitbc.bubuit.net/api/exchange_mock_api.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
