# systemd/aitbc-marketplace-enhanced.service
[Unit]
Description=AITBC Enhanced Marketplace Service (Port 8014)
Documentation=https://docs.aitbc.bubuit.net
After=network.target aitbc-coordinator-api.service postgresql.service
Wants=aitbc-coordinator-api.service postgresql.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/coordinator-api
Environment=PYTHONPATH=/opt/aitbc/apps/coordinator-api/src
Environment=PORT=8014
Environment=SERVICE_TYPE=marketplace-enhanced
Environment=DATABASE_URL=postgresql://aitbc:password@localhost:5432/aitbc
Environment=ROYALTY_ENABLED=true
Environment=LICENSING_ENABLED=true
Environment=LOG_LEVEL=INFO
ExecStart=/opt/aitbc/apps/coordinator-api/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8014
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-marketplace-enhanced

[Install]
WantedBy=multi-user.target
