# systemd/aitbc-cross-chain-reputation.service
[Unit]
Description=AITBC Cross Chain Reputation Service
After=network.target

[Service]
Type=simple
User=aitbc
WorkingDirectory=/opt/aitbc/services
ExecStart=/usr/bin/python3 -m cross_chain_reputation
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
