# systemd/aitbc-exchange-frontend.service
[Unit]
Description=AITBC Exchange Frontend Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/trade-exchange
Environment=PATH=/opt/aitbc/apps/coordinator-api/.venv/bin
Environment=PYTHONPATH=/opt/aitbc/apps/trade-exchange
ExecStart=/opt/aitbc/apps/coordinator-api/.venv/bin/python server.py --port 3002
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-exchange-frontend

[Install]
WantedBy=multi-user.target
