# systemd/aitbc-gpu-multimodal.service
[Unit]
Description=AITBC GPU Multi-Modal Processing Service
After=network.target aitbc-coordinator-api.service
Wants=aitbc-coordinator-api.service

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/home/oib/aitbc/apps/coordinator-api
Environment=PATH=/usr/bin
Environment=CUDA_VISIBLE_DEVICES=0
ExecStart=/usr/bin/python3 -m uvicorn src.app.services.gpu_multimodal_app:app --host 127.0.0.1 --port 8003
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-gpu-multimodal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/oib/aitbc/apps/coordinator-api

# GPU Access
DeviceAllow=/dev/nvidia0 rwm
DevicePolicy=auto

[Install]
WantedBy=multi-user.target
