# systemd/aitbc-coordinator-api.service
[Unit]
Description=AITBC Coordinator API Service (Python 3.13.5+)
After=network.target
Wants=network.target

[Service]
Type=simple
User=oib
Group=oib
WorkingDirectory=/opt/coordinator-api
Environment=PATH=/opt/coordinator-api/.venv/bin
Environment=PYTHONPATH=/opt/coordinator-api/src
# Python version validation
ExecStartPre=/bin/bash -c "python3 --version || (echo 'Python 3.13.5+ required' && exit 1)"
ExecStart=/opt/coordinator-api/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-coordinator-api

[Install]
WantedBy=multi-user.target
