# systemd/aitbc-adaptive-learning.service
[Unit]
Description=AITBC Adaptive Learning Service (Port 8013)
Documentation=https://docs.aitbc.bubuit.net
After=network.target aitbc-coordinator-api.service
Wants=aitbc-coordinator-api.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/coordinator-api
Environment=PYTHONPATH=/opt/aitbc/apps/coordinator-api/src
Environment=PORT=8013
Environment=SERVICE_TYPE=adaptive-learning
Environment=LEARNING_MODE=online
Environment=LOG_LEVEL=INFO
ExecStart=/opt/aitbc/.venv/bin/python -m aitbc_adaptive_learning.main
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-adaptive-learning

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aitbc/logs /opt/aitbc/data /opt/aitbc/models
LimitNOFILE=65536

# Resource limits
MemoryMax=3G
CPUQuota=250%

[Install]
WantedBy=multi-user.target
