# systemd/aitbc-marketplace.service
[Unit]
Description=AITBC Marketplace UI
After=network.target

[Service]
Type=exec
User=root
WorkingDirectory=/opt/aitbc/apps/marketplace-ui
Environment=PATH=/opt/aitbc/.venv/bin
ExecStart=/opt/aitbc/.venv/bin/python server.py --port 3001
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
