# systemd/aitbc-multimodal.service
[Unit]
Description=AITBC Multi-Modal Agent Service
After=network.target aitbc-coordinator-api.service
Wants=aitbc-coordinator-api.service

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/home/oib/aitbc/apps/coordinator-api
Environment=PATH=/usr/bin
ExecStart=/usr/bin/python3 -m uvicorn src.app.services.multimodal_app:app --host 127.0.0.1 --port 8020
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-multimodal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/oib/aitbc/apps/coordinator-api

[Install]
WantedBy=multi-user.target
