# systemd/aitbc-agent-registry.service
[Unit]
Description=AITBC Agent Registry Service
After=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/agent-services/agent-registry/src
Environment=PYTHONPATH=/opt/aitbc
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
