# systemd/aitbc-explorer.service
[Unit]
Description=AITBC Blockchain Explorer UI
After=network.target aitbc-blockchain-rpc.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/blockchain-explorer
Environment=PATH=/opt/aitbc/apps/coordinator-api/.venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/aitbc/apps/blockchain-explorer
ExecStart=/opt/aitbc/apps/coordinator-api/.venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-explorer

[Install]
WantedBy=multi-user.target
