# systemd/aitbc-agent-coordinator.service
[Unit]
Description=AITBC Agent Coordinator Service
After=network.target aitbc-agent-registry.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/agent-services/agent-coordinator/src
Environment=PYTHONPATH=/opt/aitbc
ExecStart=/usr/bin/python3 coordinator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
