# systemd/aitbc-explorer.service
[Unit]
Description=AITBC Blockchain Explorer UI
After=network.target aitbc-blockchain-rpc-1.service aitbc-blockchain-rpc-2.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/blockchain-explorer
# Using the blockchain node venv since the coordinator one is broken
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/aitbc/.env
ExecStart=/opt/aitbc/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-explorer

# Allow database writes for SQLite WAL mode
ProtectSystem=no

[Install]
WantedBy=multi-user.target
