# systemd/aitbc-ai-service.service
[Unit]
Description=AITBC AI Service
After=network.target aitbc-agent-registry.service aitbc-agent-coordinator.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/ai-engine/src
Environment=PYTHONPATH=/opt/aitbc
ExecStart=/usr/bin/python3 ai_service.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
