# systemd/aitbc-gpu-registry.service
[Unit]
Description=AITBC GPU Registry Demo
After=network.target
Wants=network.target

[Service]
Type=simple
User=oib
Group=oib
WorkingDirectory=/opt/gpu-miner
ExecStart=/usr/bin/python3 /opt/gpu-miner/gpu_registry_demo.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-gpu-registry

[Install]
WantedBy=multi-user.target
