# systemd/aitbc-gpu-registry.service
[Unit]
Description=AITBC GPU Registry Demo
After=network.target
Wants=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/gpu-miner
Environment=PATH=/usr/bin
Environment=PYTHONPATH=/opt/aitbc/apps/gpu-miner
ExecStart=/usr/bin/python3 /opt/aitbc/dev/gpu/gpu_registry_demo.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-gpu-registry

[Install]
WantedBy=multi-user.target
