# systemd/aitbc-agent-registry.service
[Unit]
Description=AITBC Agent Registry Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc/apps/agent-services/agent-registry/src
Environment=PYTHONPATH=/opt/aitbc
ExecStart=/opt/aitbc/venv/bin/python app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
