# systemd/aitbc-gpu.service
[Unit]
Description=AITBC GPU Service
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/gpu-service
Environment="PATH=/opt/aitbc/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/aitbc/apps/gpu-service/src:/opt/aitbc/packages/py/aitbc-core/src:$PYTHONPATH"
ExecStart=/opt/aitbc/venv/bin/python -m gpu_service.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-gpu

[Install]
WantedBy=multi-user.target
