# systemd/aitbc-wallet.service
[Unit]
Description=AITBC Wallet Daemon Service (Python 3.13.5+)
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/wallet-daemon
Environment=PATH=/opt/wallet-daemon/.venv/bin
Environment=PYTHONPATH=/opt/wallet-daemon
# Python version validation
ExecStartPre=/bin/bash -c "python3 --version || (echo 'Python 3.13.5+ required' && exit 1)"
ExecStart=/opt/wallet-daemon/.venv/bin/python -m wallet_daemon.main --host 0.0.0.0 --port 8001
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-wallet

[Install]
WantedBy=multi-user.target
