# systemd/aitbc-wallet.service
[Unit]
Description=AITBC Wallet Daemon Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/wallet-daemon/src
Environment=PYTHONPATH=/opt/aitbc/packages/py/aitbc-crypto/src:/opt/aitbc/packages/py/aitbc-sdk/src
ExecStart=/opt/aitbc/apps/wallet-daemon/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8002
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
