# systemd/aitbc-agent-registry.service
[Unit]
Description=AITBC Agent Registry Service
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
EnvironmentFile=/etc/aitbc/.env
EnvironmentFile=/etc/aitbc/node.env
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/scripts/wrappers/aitbc-agent-registry-wrapper.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
