# systemd/aitbc-load-secrets.service
[Unit]
Description=Load AITBC secrets at boot
Before=aitbc-blockchain-rpc.service aitbc-coordinator-api.service aitbc-wallet.service aitbc-agent-coordinator.service
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/aitbc/scripts/utils/load-keystore-secrets.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
