# systemd/aitbc-openclaw-ai.service
[Unit]
Description=AITBC OpenClaw AI Service
After=network.target aitbc-mining-blockchain.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ID=aitbc
Environment=PYTHONPATH=/opt/aitbc/production/services
EnvironmentFile=/opt/aitbc/production/.env

# OpenClaw AI execution
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/production/services/openclaw_ai.py
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10

# AI service reliability
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# AI logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-openclaw-ai

# AI security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/aitbc/production/data/openclaw /opt/aitbc/production/logs/openclaw

# AI performance
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=60%

[Install]
WantedBy=multi-user.target
