# systemd/aitbc-exchange-api.service
[Unit]
Description=AITBC Exchange API Service (Python 3.13.5+)
After=network.target
Wants=network.target

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc/apps/trade-exchange
Environment=PATH=/opt/aitbc/apps/coordinator-api/.venv/bin
# Python version validation
ExecStartPre=/bin/bash -c "python3 --version || (echo 'Python 3.13.5+ required' && exit 1)"
ExecStart=/opt/aitbc/apps/coordinator-api/.venv/bin/python simple_exchange_api.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-exchange-api

[Install]
WantedBy=multi-user.target
