# systemd/aitbc-production-miner.service
[Unit]
Description=AITBC Production GPU Miner
After=network.target ollama.service aitbc-coordinator-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/miner
Environment="PATH=/opt/aitbc/venv/bin"
Environment="PYTHONPATH=/opt/aitbc"
Environment="COORDINATOR_URL=http://localhost:8011"
Environment="MINER_ID=aitbc-miner-1"
Environment="AUTH_TOKEN=aitbc-miner-token-secure"
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/miner/production_miner.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
