# systemd/aitbc-openclaw.service
[Unit]
Description=AITBC OpenClaw Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/aitbc/apps/openclaw-service
Environment="PYTHONPATH=/opt/aitbc/apps/openclaw-service/src:/opt/aitbc"
ExecStart=/root/.local/bin/poetry run python -m openclaw_service.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
