# systemd/aitbc-loadbalancer-geo.service
[Unit]
Description=AITBC Geographic Load Balancer (Port 8017)
After=network.target aitbc-coordinator-api.service aitbc-marketplace-enhanced.service
Wants=aitbc-coordinator-api.service aitbc-marketplace-enhanced.service

[Service]
Type=simple
User=aitbc
Group=aitbc
WorkingDirectory=/opt/aitbc
Environment=PATH=/usr/bin
Environment=PORT=8017
Environment=SERVICE_TYPE=loadbalancer-geo
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 /opt/aitbc/apps/coordinator-api/scripts/geo_load_balancer.py --port 8017
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-loadbalancer-geo

[Install]
WantedBy=multi-user.target
