# systemd/aitbc-advanced-learning.service
[Unit]
Description=AITBC Advanced Learning Service
After=network.target

[Service]
Type=simple
User=aitbc
WorkingDirectory=/opt/aitbc/services
ExecStart=/usr/bin/python3 -m advanced_learning
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
