# systemd/aitbc-exchange-api.service.d/override.conf
[Service]
ExecStart=
ExecStart=/opt/aitbc/apps/coordinator-api/.venv/bin/python simple_exchange_api.py --port 8001
