# systemd/aitbc-blockchain-sync.service.d/10-central-env.conf
[Service]
EnvironmentFile=/etc/aitbc/.env
