# apps/agent-services/aitbc-agent-daemon.service
[Unit]
Description=AITBC Autonomous Agent Listener Daemon
Documentation=https://github.com/aitbc/blockchain
After=network.target aitbc-blockchain-node.service
Requires=aitbc-blockchain-node.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/node.env
Environment="PATH=/opt/aitbc/venv/bin:/usr/local/bin:/usr/bin:/bin"
# Hermes polling configuration (optional)
Environment="ENABLE_HERMES_POLLING=true"
Environment="HERMES_AGENT_IDS=owl-aitbc3"
Environment="HERMES_COORDINATOR_URL=http://hub.aitbc.bubuit.net:8011"
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/agent-services/aitbc-agent-daemon-wrapper.py

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=AgentDaemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
