# apps/wallet/aitbc-wallet.service
[Unit]
Description=AITBC Wallet Daemon with Multi-Chain Support
After=network.target aitbc-blockchain-node.service
Wants=aitbc-blockchain-node.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aitbc
EnvironmentFile=/etc/aitbc/blockchain.env
EnvironmentFile=/etc/aitbc/blockchain-secrets.env
EnvironmentFile=/etc/aitbc/node.env
Environment=WALLET_IMPORT_PASSWORD=Aitbc-Password-123
ExecStart=/opt/aitbc/venv/bin/python /opt/aitbc/apps/wallet/aitbc-wallet-wrapper.py
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=false
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aitbc-wallet

# Allow database writes for SQLite WAL mode
ProtectSystem=no

[Install]
WantedBy=multi-user.target
